# system_files/etc/profile.d/wjust.sh
#!/bin/bash
# wjust alias
alias wjust='just --justfile /usr/share/wepos/justfile'
